# external/board/rg35xxsp/overlay/etc/init.d/wifi
#!/sbin/openrc-run
#
# Connect the seeded Wi-Fi profile and provide OpenRC's virtual net service.
# Stock wpa_supplicant owns the daemon; this helper only performs board-specific
# provisioning, association, and DHCP bring-up.
#

description="Bring up Wi-Fi and provide net"

wifi_interface="wlan0"
wifi_driver_module="rtw88_8821cs"
target_ssid_file="/mnt/sdcard/.sp/config/wifi/target-ssid"
wpa_supplicant_state_dir="/etc/wpa_supplicant"
wpa_supplicant_config_file="${wpa_supplicant_state_dir}/wpa_supplicant.conf"
wifi_driver_wait_seconds=3
wpa_cli_wait_seconds=20
station_connect_wait_seconds=12
ipv4_wait_seconds=8
udhcpc_pidfile="/run/udhcpc.${wifi_interface}.pid"
udhcpc_logfile="/run/udhcpc.${wifi_interface}.log"
diagnostic_logfile="/tmp/wifi.diagnostics"

wpa_supplicant_seed_dir="/mnt/sdcard/.sp/config/wifi/wpa_supplicant"
wifi_config_file="/mnt/sdcard/.sp/config/wifi/wifi.config"

depend() {
  need localmount dbus wpa_supplicant
  after localmount wpa_supplicant
  provide net
}

# Create a default wpa_supplicant.conf if it does not exist
init_wpa_supplicant_conf() {
  mkdir -p "${wpa_supplicant_state_dir}"
  chmod 755 "${wpa_supplicant_state_dir}" 2>/dev/null || true
  
  {
    printf '%s\n' 'ctrl_interface=DIR=/var/run/wpa_supplicant GROUP=wheel'
    printf '%s\n' 'update_config=1'
    printf '%s\n' ''
  } > "${wpa_supplicant_config_file}"
  chmod 600 "${wpa_supplicant_config_file}" 2>/dev/null || true
}

# Seed profiles from the SD card wpa_supplicant/ directory
seed_wpa_profiles() {
  [ -d "${wpa_supplicant_seed_dir}" ] || return 0

  for profile_path in "${wpa_supplicant_seed_dir}"/*.conf; do
    [ -f "${profile_path}" ] || continue
    cat "${profile_path}" >> "${wpa_supplicant_config_file}"
    printf '\n' >> "${wpa_supplicant_config_file}"
  done
}

# Parse wifi.config SSID and Passphrase and generate a network block
generate_wpa_profile_from_config() {
  [ -f "${wifi_config_file}" ] || return 0

  config_ssid=""
  config_pass=""

  while IFS='=' read -r key val; do
    [ -n "$key" ] || continue
    case "$key" in
      \#*) continue ;;
      SSID) config_ssid="$val" ;;
      Passphrase) config_pass="$val" ;;
    esac
  done < "${wifi_config_file}"

  [ -n "$config_ssid" ] || return 0
  [ -n "$config_pass" ] || return 0

  {
    printf '%s\n' 'network={'
    printf '\t%s\n' "ssid=\"$config_ssid\""
    printf '\t%s\n' "psk=\"$config_pass\""
    printf '%s\n' '}'
    printf '%s\n' ''
  } >> "${wpa_supplicant_config_file}"
}

prepare_wifi_interface() {
  i=0

  if command -v rfkill >/dev/null 2>&1; then
    rfkill unblock wifi >/dev/null 2>&1 || true
  fi



  while [ "${i}" -lt "${wifi_driver_wait_seconds}" ]; do
    if ip link show "${wifi_interface}" >/dev/null 2>&1; then
      ip link set "${wifi_interface}" up >/dev/null 2>&1 || true
      return 0
    fi
    i=$((i + 1))
    sleep 1
  done

  return 1
}

wait_for_wpa_supplicant_ready() {
  i=0

  while [ "${i}" -lt "${wpa_cli_wait_seconds}" ]; do
    if pidof wpa_supplicant >/dev/null 2>&1 &&
      wpa_cli -i "${wifi_interface}" ping >/dev/null 2>&1; then
      return 0
    fi
    i=$((i + 1))
    sleep 1
  done

  return 1
}

log_failure_diagnostics() {
  reason="$1"

  {
    printf '%s\n' "reason=${reason}"
    printf '%s ' "/sys/class/net/${wifi_interface}:"
    if [ -e "/sys/class/net/${wifi_interface}" ]; then
      printf '%s\n' "present"
      ip link show "${wifi_interface}" 2>&1 || true
    else
      printf '%s\n' "missing"
    fi

    printf '%s\n' "/sys/bus/sdio/devices:"
    ls -la /sys/bus/sdio/devices 2>&1 || true

    printf '%s\n' "wpa_cli -i ${wifi_interface} status:"
    if command -v wpa_cli >/dev/null 2>&1; then
      wpa_cli -i "${wifi_interface}" status 2>&1 || true
    else
      printf '%s\n' "wpa_cli missing"
    fi
  } >> "${diagnostic_logfile}"

  logger -t wifi "startup failed: ${reason}; diagnostics in ${diagnostic_logfile}" \
    2>/dev/null || true
  echo "wifi: startup failed: ${reason}" >&2
}

station_has_connection() {
  wpa_cli -i "${wifi_interface}" status 2>/dev/null | grep -q "wpa_state=COMPLETED"
}

wait_for_station_connected() {
  i=0

  while [ "${i}" -lt "${station_connect_wait_seconds}" ]; do
    if station_has_connection; then
      return 0
    fi
    i=$((i + 1))
    sleep 1
  done

  return 1
}

wait_for_ipv4_address() {
  i=0

  while [ "${i}" -lt "${ipv4_wait_seconds}" ]; do
    if ip -4 addr show "${wifi_interface}" 2>/dev/null | grep -q "inet "; then
      return 0
    fi
    i=$((i + 1))
    sleep 1
  done

  return 1
}

start_dhcp_client() {
  rm -f "${udhcpc_pidfile}"
  : > "${udhcpc_logfile}"

  udhcpc -b -R -O search -O staticroutes \
    -p "${udhcpc_pidfile}" \
    -i "${wifi_interface}" \
    -s /usr/share/udhcpc/default.script \
    >> "${udhcpc_logfile}" 2>&1
}

start() {
  init_wpa_supplicant_conf
  generate_wpa_profile_from_config
  seed_wpa_profiles

  # If no network blocks were added, there is nothing to connect to.
  # Exit gracefully so the service does not block ntpd or spam the console.
  if ! grep -q '^network={' "${wpa_supplicant_config_file}" 2>/dev/null; then
    einfo "No Wi-Fi profiles configured; skipping connection"
    return 0
  fi

  if ! prepare_wifi_interface; then
    log_failure_diagnostics "missing-${wifi_interface}"
    return 1
  fi

  if ! wait_for_wpa_supplicant_ready; then
    log_failure_diagnostics "wpa_supplicant-not-ready"
    return 1
  fi

  # Tell wpa_supplicant to re-read config if it is already running
  wpa_cli -i "${wifi_interface}" reconfigure >/dev/null 2>&1 || true

  if ! wait_for_station_connected; then
    log_failure_diagnostics "station-connect-timeout"
    return 1
  fi

  wait_for_ipv4_address && return 0
  if ! start_dhcp_client; then
    log_failure_diagnostics "dhcp-start-failed"
    return 1
  fi
  if ! wait_for_ipv4_address; then
    log_failure_diagnostics "ipv4-timeout"
    return 1
  fi
  return 0
}

stop() {
  if [ -f "${udhcpc_pidfile}" ]; then
    kill "$(cat "${udhcpc_pidfile}")" >/dev/null 2>&1 || true
    rm -f "${udhcpc_pidfile}"
  fi
  ip link set "${wifi_interface}" down >/dev/null 2>&1 || true
}
