# alpine/aports/minime-overlay/files/etc/init.d/panfrost
#!/sbin/openrc-run
# shellcheck disable=SC2034
# panfrost: load the Panfrost kernel module for the GLES path used
# by the Alpine userspace (Mesa).

description="Load Panfrost GPU driver (Alpine GLES provider)"

depend() {
	need modules
	before ui dbus
}

start() {
	ebegin "Loading Panfrost"
	modprobe panfrost 2>/dev/null || true
	eend 0
}

stop() {
	rmmod panfrost 2>/dev/null || true
	return 0
}
