# alpine/aports/minime-overlay/files/etc/init.d/telnetd
#!/sbin/openrc-run
# shellcheck disable=SC2034
# telnetd: passwordless telnet on TCP 23, autologin shell.

description="Minime telnet daemon"

command="/usr/sbin/tcpsvd"
command_args="-vE 0.0.0.0 23 /usr/sbin/telnetd -i -l /usr/bin/autologin"
pidfile="/run/telnetd.pid"

depend() {
	need modules
}

start() {
	ebegin "Starting telnetd"
	start-stop-daemon --start --background --make-pidfile \
		--pidfile "$pidfile" --exec "$command" -- $command_args
	eend $?
}

stop() {
	ebegin "Stopping telnetd"
	start-stop-daemon --stop --pidfile "$pidfile" --exec "$command"
	eend $?
}
