# alpine/aports/minime-overlay/files/etc/init.d/traits
#!/sbin/openrc-run
# shellcheck disable=SC2034
# traits: write /mnt/sdcard/.minime/traits from the immutable
# /usr/share/minime/traits payload.

description="Generate Minime device traits"

depend() {
	need modules
	before ui
	after fb-unblank
}

start() {
	ebegin "Generating Minime traits"
	/usr/share/minime/scripts/traits.sh start
	eend $?
}
