# alpine/alpine/aports/minime-overlay/files/etc/init.d/minime-ftpd
#!/sbin/openrc-run
# shellcheck disable=SC2034
# minime-ftpd: anonymous-feel FTP to the SD card.  Mirrors the
# Passwordless local FTP service.

description="Minime FTP daemon"

command="/usr/sbin/tcpsvd"
command_args="-vE 0.0.0.0 21 /usr/sbin/ftpd -w /mnt/sdcard"
pidfile="/run/ftpd.pid"

depend() {
	need minime-modules
}

start() {
	ebegin "Starting ftpd"
	start-stop-daemon --start --background --make-pidfile \
		--pidfile "$pidfile" --exec "$command" -- $command_args
	eend $?
}

stop() {
	ebegin "Stopping ftpd"
	start-stop-daemon --stop --pidfile "$pidfile" --exec "$command"
	eend $?
}
