# alpine/alpine/aports/minime-overlay/files/etc/init.d/minime-wifi
#!/sbin/openrc-run
# shellcheck disable=SC2034
# minime-wifi: wpa_supplicant + udhcpc. Keep behavior aligned with
# /usr/share/minime/scripts/minime-wifi.

description="Minime Wi-Fi (wpa_supplicant + udhcpc)"

command=":"
command_background="yes"
pidfile="/run/minime-wifi.pid"

depend() {
	need minime-modules
	after minime-fb-unblank
}

start() {
	ebegin "Starting Wi-Fi"
	/usr/share/minime/scripts/minime-wifi start &
	eend 0
}

stop() {
	ebegin "Stopping Wi-Fi"
	/usr/share/minime/scripts/minime-wifi stop
	eend $?
}
