# alpine/board/rg35xxsp/overlay/etc/init.d/telnetd
#!/sbin/openrc-run

description="BusyBox Telnet daemon"

command="/usr/sbin/telnetd"
command_args="-F -l /usr/bin/autologin"
command_background="yes"
pidfile="/run/telnetd.pid"

depend() {
	need localmount
	after bootmisc
}
