# board/rg35xxsp/overlay/etc/init.d/ftpd
#!/sbin/openrc-run

description="BusyBox FTP daemon"

command="/usr/bin/tcpsvd"
command_args="-vE 0.0.0.0 21 /usr/sbin/ftpd -w /"
command_background="yes"
pidfile="/run/ftpd.pid"

depend() {
	need localmount
	after bootmisc
}
