# alpine/external/board/common/overlay/etc/inittab
# /etc/inittab for minime handheld
# No getty on tty1 (framebuffer) to prevent login prompt from overlaying MinUI.
# Serial console on ttyS0 is available for debugging.

::sysinit:/bin/mkdir -p /dev/pts /dev/shm
::sysinit:/bin/mount -a
::sysinit:/etc/init.d/rcS

ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100

::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
::restart:/sbin/init
